# etc/systemd/system/tele-entry.service
[Unit]
Description=Tele-entry

[Service]
TimeoutStartSec=1
ExecStart=/opt/tele-entry -addr "0.0.0.0:6921"

[Install]
WantedBy=multi-user.target

